# home/.config/systemd/user/kodi-backup-cleanup.timer
[Unit]
Description=regularly clean kodi backups

[Timer]
OnCalendar=Thursday *-*-* 02:00:00

[Install]
WantedBy=default.target
