# home/.config/systemd/user/redeem-points.timer
[Unit]
Description=schedule to redeem iptorrents points

[Timer]
OnCalendar=10:00
