# home/.config/systemd/user/redeem-points.service
[Unit]
Description=redeem ipt bonus points

[Service]
Type=oneshot
ExecStart=/home/tylerjl/.bin/ribp.rb
StandardOutput=journal
