# home/.config/systemd/user/torrent-prune.timer
[Unit]
Description=regularly prune extra torrents

[Timer]
OnCalendar=11:00

[Install]
WantedBy=default.target
