# home/.config/systemd/user/zbc.timer
[Unit]
Description=schedule to check ZFS backup replication

[Timer]
OnCalendar=12:00
