# home/.config/systemd/user/zfs-scrub.timer
[Unit]
Description=schedule regular zfs scrubs

[Timer]
OnCalendar=Monday *-*-* 04:00:00

[Install]
WantedBy=default.target
