# home/.config/systemd/user/zfs-check.service
[Unit]
Description=check whether the zpool is healthy

[Service]
Type=oneshot
ExecStart=/home/tylerjl/.bin/zfs-check.rb
