# home/.config/systemd/user/zfs-scrub.timer
[Unit]
Description=schedule regular zfs scrubs

[Timer]
OnCalendar=Sunday *-*-* 02:00:00
