# home/.config/systemd/user/zfs-check.timer
[Unit]
Description=schedule checking whether the zpool is healthy

[Timer]
OnCalendar=Monday *-*-* 12:00:00
