# home/.config/systemd/user/zfs-scrub.service
[Unit]
Description=perform a zfs scrub

[Service]
Type=oneshot
ExecStart=/usr/bin/sudo zpool scrub tank
